# rc.local
#To increase track point sensitivity
echo 250 > /sys/bus/serio/devices/serio2/sensitivity
